# etc/systemd/system/transmission.service.d/transmission.conf
[Unit]
RequiresMountsFor=/mnt/sda1/transmission/done
